# monitor/debian/pico-monitor.service
[Unit]
Description=Pico LCD 系统硬件监控服务
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/pico-monitor --config /etc/pico-monitor.conf
Restart=always
RestartSec=3
User=root
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict

[Install]
WantedBy=multi-user.target
